# srcs/wordpress/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm7.3 --nodaemonize --allow-to-run-as-root
autorestart=false
startsecs = 0